# packaging/kcd-system.service
# System-level kcd service template.
# Runs the daemon as a named user, managed by root / a system administrator.
#
# Enable for a specific user:
#   sudo systemctl enable --now kcd@alice
#
# For most desktop setups, the USER service (kcd-user.service) is preferred.
# Use this system unit for headless servers or when you want kcd to start
# before any user session (e.g. to keep the phone connected at all times).

[Unit]
Description=kcd — Headless KDE Connect Daemon (%i)
Documentation=https://github.com/bethropolis/kcd
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main

# %i is the instance name — the username passed to `systemctl enable kcd@<user>`.
User=%i
Group=%i

ExecStart=/usr/bin/kcd daemon

KillMode=mixed
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

# ── XDG paths ─────────────────────────────────────────────────────────────────
# When running as a system service there is no login session, so XDG variables
# must be set explicitly. Adjust UID to match the target user if needed.
Environment="HOME=/home/%i"
Environment="XDG_CONFIG_HOME=/home/%i/.config"
Environment="XDG_STATE_HOME=/home/%i/.local/state"
Environment="XDG_CACHE_HOME=/home/%i/.cache"
# XDG_RUNTIME_DIR is user-session-specific; the daemon creates it if absent.
Environment="XDG_RUNTIME_DIR=/run/user/%U"

# ── Filesystem hardening ──────────────────────────────────────────────────────
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/.config/kcd /home/%i/.local/state/kcd /home/%i/Downloads/kcd

RuntimeDirectory=kcd
RuntimeDirectoryMode=0700

PrivateTmp=true
NoNewPrivileges=true
MemoryMax=128M
TasksMax=64
PrivateNetwork=false

[Install]
WantedBy=multi-user.target
